# daemon/com.joshua.arthur.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joshua.arthur</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/joshua/Documents/Code/arthur/daemon/arthur_watchdog.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/joshua/Documents/Code/arthur/logs/daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/joshua/Documents/Code/arthur/logs/daemon_error.log</string>
    
    <key>KeepAlive</key>
    <false/>

    <key>Disabled</key>
    <true/>
    
    <!-- Restart every hour if it crashes -->
    <key>ThrottleInterval</key>
    <integer>3600</integer>
    
    <!-- Run with background priority to avoid starving the system -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
